# cron.d/battery-mounter
# /etc/cron.d/battery-mounter: crontab entries for the battery-mounter package

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* *    * * *   root	battery-mounter.py >/dev/null
